# conf/docker-concrete_mixer.service
[Unit]
Description=Concrete Mixer service

[Service]
Restart=always
ExecStart=/concrete-mixer/init.sh
User=cmuser

[Install]
WantedBy=default.target
